# crontab
* * * * * python3 store_locations.py; echo 'ran crontab'
* * * * * sleep 30; python3 store_locations.py

